# LaunchDaemons/com.apple.awdd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>EnablePressuredExit</key>
	<true/>
	<key>Label</key>
	<string>com.apple.awdd</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.notifyd.matching</key>
		<dict>
			<key>com.apple.awdd.ConfigChange</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.awdd.ConfigChange</string>
			</dict>
			<key>com.apple.awdd.AnonymityChangedCommandLine</key>
			<dict>
				<key>Notification</key>
				<string>com.apple.awdd.anonymity</string>
			</dict>
		</dict>
	</dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.awdd</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/WirelessDiagnostics.framework/Support/awdd</string>
	</array>
</dict>
</plist>
