# LaunchDaemons/com.apple.cmio.AVCAssistant.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.cmio.AVCAssistant</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.cmio.AVCAssistant</key>
		<true/>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/Frameworks/CoreMediaIO.framework/Resources/AVC.plugin/Contents/Resources/AVCAssistant</string>
	</array>
	<key>ThrottleInterval</key>
	<integer>0</integer>
	<key>EnableTransactions</key>
	<true/>
	<key>POSIXSpawnType</key>
	<string>Interactive</string>
	<key>UserName</key>
	<string>_cmiodalassistants</string>
	<key>GroupName</key>
	<string>_cmiodalassistants</string>
</dict>
</plist>
