# LaunchDaemons/com.apple.msrpc.lsarpc.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.msrpc.lsarpc</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/libexec/rpcsvchost</string>
		<string>-launchd</string>
		<string>-sandbox</string>
		<string>lsarpc.bundle</string>
		<string>dssetup.bundle</string>
	</array>
	<key>EnableTransactions</key>
	<false/>
	<key>Sockets</key>
	<dict>
		<key>ncalrpc</key>
		<dict>
			<key>Bonjour</key>
			<false/>
			<key>SockPathMode</key>
			<integer>438</integer>
			<key>SockPathName</key>
			<string>/var/rpc/ncalrpc/lsarpc</string>
		</dict>
		<key>ncacn_np</key>
		<dict>
			<key>Bonjour</key>
			<false/>
			<key>SockPathMode</key>
			<integer>438</integer>
			<key>SockPathName</key>
			<string>/var/rpc/ncacn_np/lsarpc</string>
		</dict>
	</dict>
	<key>EnvironmentVariables</key>
	<dict>
		<key>RPC_DEBUG</key>
		<string>0</string>
	</dict>
</dict>
</plist>
