# LaunchAgents/com.apple.corespeechd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>MachServices</key>
	<dict>
		<key>com.apple.corespeech.corespeechservices</key>
		<true/>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>Label</key>
	<string>com.apple.corespeechd</string>
	<key>ProcessType</key>
	<string>Interactive</string>
	<key>EnableTransactions</key>
	<true/>
	<key>EnablePressuredExit</key>
	<true/>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/PrivateFrameworks/CoreSpeech.framework/corespeechd</string>
	</array>
	<key>KeepAlive</key>
	<true/>
</dict>
</plist>
