# LaunchAgents/com.apple.iCloudUserNotifications.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.iCloudUserNotificationsd</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.icloudusernotification.agent</key>
		<true/>
	</dict>
	<key>Program</key>
	<string>/System/Library/PrivateFrameworks/AOSAccounts.framework/Versions/A/Resources/iCloudUserNotificationsd.app/Contents/MacOS/iCloudUserNotificationsd</string>
	<key>EnableTransactions</key>
	<true/>
	<key>LimitLoadToSessionType</key>
	<array>
		<string>LoginWindow</string>
		<string>Aqua</string>
	</array>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.usernotificationcenter.matching</key>
		<dict>
			<key>iCloudUserNotification</key>
			<dict>
				<key>bundleid</key>
				<string>com.apple.iCloudUserNotificationsd</string>
				<key>delay registration</key>
				<true/>
				<key>events</key>
				<array>
					<string>didDismissAlert</string>
					<string>didActivateNotification</string>
					<string>didDeliverNotification</string>
					<string>didSnoozeAlert</string>
					<string>didRemoveDeliveredNotifications</string>
					<string>didExpireNotifications</string>
				</array>
			</dict>
		</dict>
	</dict>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ExitTimeOut</key>
	<integer>1</integer>
</dict>
</plist>
